# systemd/subwave-web.service
[Unit]
Description=SUB/WAVE — Next.js web UI
PartOf=subwave.target

[Service]
Type=simple
# Next.js standalone output bundles its own minimal server.js + the runtime
# files into .next/standalone. install-native.sh runs `npm ci && npm run build`
# in web/ to produce it. Static assets (.next/static, public/) live next to it
# only after the installer copies them — see install-native.sh.
WorkingDirectory=${REPO}/web/.next/standalone

# Production mode. The API/stream URLs are baked into the build by
# install-native.sh: same-origin defaults (/api, /stream.mp3) behind the Caddy
# edge, or absolute URLs when the installer is run with SUBWAVE_EDGE=none.
# NEXT_PUBLIC_* are build-time constants — setting them here at runtime would
# have no effect; re-run the installer to change them.
#
# HOSTNAME is rendered by the installer: 127.0.0.1 behind the Caddy edge
# (only Caddy reaches it), 0.0.0.0 under SUBWAVE_EDGE=none (browsers hit it).
Environment=NODE_ENV=production
Environment=PORT=7700
Environment=HOSTNAME=${WEB_HOSTNAME}

ExecStart=${NODE_BIN} server.js
Restart=on-failure
RestartSec=3s

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=subwave.target
